# go-present.service
[Unit]
Description = sss server daemon

[Service]
ExecStart = /opt/sss
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target


